# cti-newsfeed.timer
[Unit]
Description=CTI News Feed — Gunluk zamanlayici (09:00)

[Timer]
# Her gun Istanbul saatiyle 09:00'da calistir (sunucu UTC ise 06:00 UTC)
OnCalendar=*-*-* 06:00:00

# Sunucu kapali iken atlanan calistirmayi, acilinca otomatik yap
Persistent=true

# Birden fazla instance cakismasin
AccuracySec=60

[Install]
WantedBy=timers.target
