# cti-newsfeed.service
[Unit]
# Systemd'de gorulecek aciklama (systemctl status komutunda gosterilir)
Description=CTI News Feed Automation — Threat Briefing

# Internet baglantisi hazir olana kadar bekle (Gemini API ve SMTP icin)
After=network-online.target
Wants=network-online.target

[Service]
# oneshot = bir kez calisir ve cikar (daemon degil, gunluk gorev)
Type=oneshot

# Script bu klasorde calisacak (log dosyalari icin onemli)
WorkingDirectory=/opt/cti-project

# Python'u venv icinden cagir (sistem Python'unu kullanma)
ExecStart=/opt/cti-project/venv/bin/python3 /opt/cti-project/cti_automation.py

# Yetkisiz sistem kullanicisi (root degil) — guvenlik icin sart
User=cti
Group=cti

# Yeni olusturulan log dosyalarinin izni 640 olur (others okuyamaz)
UMask=0137

# ── Yeniden deneme politikasi ────────────────────────
# Auto-restart KAPALI. Sebebi: acik kalirsa tek bir hata systemd tarafindan
# her 5 dk'da bir sinirsiz yeniden baslatilarak Gemini API kotasini (RPD) tuketir.
# (Gecmiste: saatte 12 restart × 3 script-ici retry = 36 istek → 20 RPD limiti asildi.)
# Gecici hatalar zaten script ici retry ile karsilanir (3 deneme, 2-4-8sn backoff).
# Calisma basarisiz olursa ertesi gun timer normal sekilde yeniden tetikler.
Restart=no

# ── Kaynak limitleri ─────────────────────────────────
# Maks RAM kullanimi (RSS parse + Gemini + mail icin yeterli)
MemoryMax=512M
# Swap kullanma — surec hizli kalsin
MemorySwapMax=0
# Tek cekirdek %80'i (sunucuyu zorlamasin)
CPUQuota=80%
# Process/thread sayisi limiti
TasksMax=32
# 10 dakikada bitmezse durdur (normalde ~1 dk surer)
TimeoutStartSec=600

# ── Dosya sistemi izolasyonu ─────────────────────────
# Tum dosya sistemini sadece okunabilir yap
ProtectSystem=strict
# /home klasorlerini gizle (script home'a erisemez)
ProtectHome=true
# Sadece bu klasore yazabilir (log dosyalari icin)
ReadWritePaths=/opt/cti-project/logs
# Her calistirmada ayri /tmp ortami (paylasimli /tmp riskini onler)
PrivateTmp=true
# /dev/* cihaz dosyalarina erisimi engelle
PrivateDevices=true

# ── Kernel korumalari ────────────────────────────────
# /proc/sys ve /sys icindeki kernel parametrelerini degistiremez
ProtectKernelTunables=true
# Yeni kernel modulu yukleyemez
ProtectKernelModules=true
# Kernel loglarini okuyamaz (dmesg engellenir)
ProtectKernelLogs=true
# cgroup'lari degistiremez (kaynak limit kacirma engellenir)
ProtectControlGroups=true
# Sistem saatini degistiremez
ProtectClock=true
# Hostname degistiremez
ProtectHostname=true
# Personality syscall'i kisitla (ASLR bypass engellenir)
LockPersonality=true

# ── Ayricalik kisitlamalari ──────────────────────────
# setuid binary calistirsa bile yetki yukseltemez
NoNewPrivileges=true
# SUID/SGID bit'li dosya olusturamaz
RestrictSUIDSGID=true
# Real-time scheduling kullanamaz
RestrictRealtime=true
# Yeni namespace olusturamaz (container kacis koruması)
RestrictNamespaces=true

# ── Ag kisitlamalari ─────────────────────────────────
# Sadece IPv4, IPv6 ve Unix soketlerine izin ver
# (AF_PACKET, AF_NETLINK gibi tehlikeli soketleri engelle)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# ── Process gorunurlugu ──────────────────────────────
# Baska kullanicilarin processlerini goremez
ProtectProc=invisible
# /proc'ta sadece process bilgilerini gor (sistem bilgisi gizli)
ProcSubset=pid

# ── Ortam degiskenleri ───────────────────────────────
# Gemini API key vb. .env dosyasindan okuniyor (script icinde dotenv ile)
# Sadece PATH lazim
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
# Multi-user mode'da etkin olsun (normal sistem baslangici)
WantedBy=multi-user.target
